# charge_limit.path
[Path]
PathModified=/etc/charge_limit

[Install]
WantedBy=multi-user.target
