# fan_profile.path
[Path]
PathModified=/etc/fan_profile

[Install]
WantedBy=multi-user.target
